# debian/database-manager.service
[Unit]
Description=Database Manager - универсальный менеджер баз данных
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/database-manager
ExecStart=/usr/bin/database-manager
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Порт настраивается через /etc/database-manager/app.json
# Для переопределения можно использовать переменную окружения:
# Environment="PORT=8081"

[Install]
WantedBy=multi-user.target

