# debian/postinst
#!/bin/bash
set -e

if [ "$1" = "configure" ]; then
    # Создаем директории для конфигурации и данных
    mkdir -p /etc/database-manager
    mkdir -p /var/lib/database-manager
    
    # Создаем начальные конфигурационные файлы, если их нет
    if [ ! -f /etc/database-manager/connections.json ]; then
        echo "[]" > /etc/database-manager/connections.json
        chmod 644 /etc/database-manager/connections.json
    fi
    
    if [ ! -f /etc/database-manager/users.json ]; then
        echo "[]" > /etc/database-manager/users.json
        chmod 644 /etc/database-manager/users.json
    fi
    
    if [ ! -f /etc/database-manager/app.json ]; then
        echo '{"port": "8081"}' > /etc/database-manager/app.json
        chmod 644 /etc/database-manager/app.json
    fi
    
    # Устанавливаем права доступа
    chown -R root:root /etc/database-manager
    chown -R root:root /var/lib/database-manager
    chown -R root:root /usr/share/database-manager
    
    # Перезагружаем systemd для обнаружения нового сервиса
    systemctl daemon-reload >/dev/null 2>&1 || true
    
    # Включаем сервис, но не запускаем автоматически
    systemctl enable database-manager.service >/dev/null 2>&1 || true
    
    echo "Database Manager установлен успешно!"
    echo "Для запуска сервиса выполните: sudo systemctl start database-manager"
    echo "Порт по умолчанию: 8081 (можно изменить в /etc/database-manager/app.json)"
fi

exit 0

